# pagewatch.service
[Unit]
Description=Run script to check stock

[Service]
Type=oneshot
ExecStart=/home/john/project/mc/script

[Install]
WantedBy=default.target
